# neptune.service
[Service]
Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neptune
ExecStart=%h/.cargo/bin/neptune

[Install]
WantedBy=multi-user.target
WantedBy=graphical.target
WantedBy=graphical-session.target