# neptune.service
[Service]
ExecStart=/home/yummi/.cargo/bin/neptune
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=neptune
User=yummi
Group=yummi
SupplementaryGroups=input

[Install]
WantedBy=multi-user.target