# neptune.service
[Unit]
Description=Cool things
After=default.target

[Service]
Restart=always
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neptune
Environment="NEPTUNE_LOG=debug"
ExecStart=%h/.cargo/bin/neptune daemon
KillSignal=SIGINT

[Install]
WantedBy=default.target
